# agentserver/data/supervisor/agentserver.conf
[program:agentserver]
command=/home/vagrant/.virtualenvs/as/bin/python /home/vagrant/src/agentserver/agentserver/agentserver.py --config=/etc/agentserver/agentserver.conf
directory=/home/vagrant/src/agentserver/agentserver
user=vagrant
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/var/log/supervisor_agentserver.log